# launchd/ai.signet.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.signet.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jason/.bun/bin/bun</string>
        <string>/Users/jason/.bun/install/global/node_modules/signetai/dist/daemon.js</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>SIGNET_PORT</key>
        <string>3850</string>
        <key>SIGNET_PATH</key>
        <string>/Users/jason/.local/share/signet</string>
        <key>PI_CODING_AGENT_DIR</key>
        <string>/Users/jason/.config/pi/agent</string>
        <key>PATH</key>
        <string>/Users/jason/.bun/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/jason/.local/share/signet/.daemon/logs/daemon.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/jason/.local/share/signet/.daemon/logs/daemon.err.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/jason/.local/share/signet</string>
</dict>
</plist>
